# scripts/systemd/ewm-one.socket
# An Apple 1 per telnet connection: inetd-style socket activation —
# systemd accepts each connection on port 6502 and starts a fresh
# ewm-one@.service instance with the connection as stdin/stdout.
#
# Build, install, and operating guide: README.md in this directory.

[Unit]
Description=Apple 1 telnet listener (one machine per connection)

[Socket]
ListenStream=6502
Accept=yes
MaxConnectionsPerSource=4

[Install]
WantedBy=sockets.target
